# ansible-jupyter/roles/jupyter/templates/jupyter.service.j2
# Managed by Ansible — do not edit manually
[Unit]
Description=Jupyter Notebook Server
After=network.target

[Service]
Type=simple
User={{ jupyter_user }}
ExecStart={{ jupyter_venv }}/bin/jupyter notebook
WorkingDirectory={{ jupyter_home }}
Restart=always

[Install]
WantedBy=multi-user.target
